# rootfs/usr/lib/systemd/system/report-updates.path
[Unit]
Description=Report new rpm-ostree deployments

[Path]
PathExists=/run/ostree/staged-deployment

[Install]
WantedBy=multi-user.target
